# layers/fuel-meter/recipes-core/fuel-meter-app/files/fuel-meter.service
[Unit]
Description=fuel-meter Logic
After=multi-user.target

[Service]
ExecStart=/usr/bin/fuel-meter-app
Restart=always

[Install]
WantedBy=multi-user.target